# supervisor/rasa-x.conf
## rasa-x.conf

[supervisord]
nodaemon=true

[program:action_server]
command=python3 -m rasa_sdk --actions actions
autostart=true
autorestart=true
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes = 0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0

[program:rasa_x]
command=python3 -m rasa x --domain=domain --rasa-x-port=5002 --no-prompt
autostart=true
autorestart=true
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes = 0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
